# contrib/leeward.user.service
[Unit]
Description=Leeward Sandbox Daemon (User)
After=default.target

[Service]
Type=simple
ExecStart=%h/.local/bin/leeward-daemon
Restart=always
RestartSec=5
RuntimeDirectory=leeward
RuntimeDirectoryMode=0700

[Install]
WantedBy=default.target