# contrib/leeward.system.service
[Unit]
Description=Leeward Sandbox Daemon
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/leeward-daemon
Restart=always
RestartSec=5
User=leeward
Group=leeward

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
RuntimeDirectory=leeward
RuntimeDirectoryMode=0755

# Resource limits
MemoryMax=1G
CPUQuota=200%

[Install]
WantedBy=multi-user.target